# supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:chat]
command=/usr/local/bin/npm start
directory=/code
autostart=true
autorestart=true
